# bash/setup-files/network-startup.service
[Unit]
After=network-online.target

[Service]
ExecStart=/usr/local/bin/network_startup.sh

[Install]
WantedBy=default.target

