# bash/setup-files/nfs-automount.service
[Unit]
After=network-online.target

[Service]
ExecStart=/usr/local/bin/nfs_automount.sh

[Install]
WantedBy=default.target

